# crontab
59 23 28-31 * * [ "$(date +\%d -d tomorrow)" == "01" ] && python /app/script.py